# AI_Voice_Assistant_using_Raspi4/Yocto/Yocto_sources/meta-userapp/recipes-connectivity/dhcpcd/files/dhcpcd.service
[Unit]
Description=DHCP Client Daemon from ferganey
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/sbin/dhcpcd -q -w wlan0
Restart=on-failure

[Install]
WantedBy=multi-user.target

